# deploy/systemd/euroqa-frontend.service
[Unit]
Description=Euro_QA Vite preview frontend
Documentation=file:/home/root251/euroqa/deploy/systemd/README.md
After=network-online.target euroqa-backend.service
Wants=network-online.target

[Service]
Type=simple
User=root251
Group=root251
WorkingDirectory=/home/root251/euroqa
Environment=PATH=/home/root251/.local/bin:/home/root251/.local/share/pnpm:/usr/local/bin:/usr/bin:/bin
ExecStartPre=/usr/bin/env pnpm --dir frontend build
ExecStart=/usr/bin/env pnpm --dir frontend preview --host 0.0.0.0 --port 4173
Restart=on-failure
RestartSec=10
TimeoutStartSec=60
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
